# move2docker/etc/service/netscript/run
#!/bin/bash
for filename in /etc/netscript/*.sh; do
    [ -e "$filename" ] || continue
    parser=$(head -n 1 "$filename")
    eval "${parser:2}" /etc/rc.local
done
sleep infinity